# service/timers/login_check.timer
[Unit]
Description=Run login_check script every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
