# service/timers/resource_usage.timer
[Unit]
Description=Run resource_usage script every 30 seconds

[Timer]
OnBootSec=30sec
OnUnitActiveSec=30sec
Persistent=true

[Install]
WantedBy=timers.target
